# roles/nfs-client/templates/mnt-nfsshare.mount.j2
[Unit]
Description = NFS share

[Mount]
What = {{ NFS_SERVER }}:/mnt/nfs
Where = /mnt/nfsshare
Type = nfs4
Options = rw
DirectoryMode = 0755